# base/etc/crontabs/nginx
* * * * * /usr/bin/php /var/www/localhost/htdocs/protected/yii queue/run
* * * * * /usr/bin/php /var/www/localhost/htdocs/protected/yii cron/run
