# etc/crontabs/nginx
* * * * * /usr/bin/php __HUMHUB_WEB_ROOT____HUMHUB_SUB_DIR__/protected/yii queue/run >/dev/null 2>&1
* * * * * /usr/bin/php __HUMHUB_WEB_ROOT____HUMHUB_SUB_DIR__/protected/yii cron/run >/dev/null 2>&1
